# deploy/gateway/cheeze-portal-api.service.example
[Unit]
Description=CHEEZE public portal control facade
After=network-online.target cheeze-control-api.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cheeze-control
Environment=CHEEZE_PORTAL_LISTEN_HOST=127.0.0.1
Environment=CHEEZE_PORTAL_LISTEN_PORT=11437
Environment=CHEEZE_INTERNAL_CONTROL_BASE=http://127.0.0.1:11436
Environment=CHEEZE_PORTAL_CONTROL_TOKEN=CHANGE_ME_WITH_A_LONG_RANDOM_SECRET
Environment=CHEEZE_PORTAL_CONTROL_HEADER=X-Cheeze-Control-Token
Environment=CHEEZE_PORTAL_REQUEST_TIMEOUT=210
Environment=CHEEZE_PORTAL_TOKEN_REGISTRY=/opt/cheeze-control/portal-control-tokens.json
Environment=CHEEZE_PORTAL_AUDIT_LOG=/opt/cheeze-control/portal-control-audit.log
ExecStart=/usr/bin/python3 /opt/cheeze-control/cheeze-portal-api.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
