# deploy/discord-bot/cheeze-discord-bot.service.example
[Unit]
Description=CHEEZE Discord game control bot
After=network-online.target cheeze-portal-api.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cheeze-bot
EnvironmentFile=/etc/cheeze-bot/cheeze-discord-bot.env
ExecStart=/usr/bin/python3 /opt/cheeze-bot/cheeze-discord-bot.py
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
